# ubuntu-zestyzapus-elastic-2-3-3/src/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
user=elastic
command=/opt/elasticsearch-2.3.3/bin/elasticsearch
autorestart=false
startretries=0


