# ubuntu-xenial-tomcat-8/src/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
user=root
command=/opt/tomcat/bin/tomcat.sh start
autorestart=false
startretries=0


