# ubuntu-zestyzapus-cassandra-212/src/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:cassandra]
user=root
command=/opt/apache-cassandra-2.1.12/bin/cassandra -f
autorestart=false
startretries=0

