# ubuntu-xenial-elastic-0-19/src/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
user=root
command=/opt/elasticsearch-0.19.9/bin/elasticsearch -f
autorestart=false
startretries=0


